# rbdbackuptool.service
[Unit]
Description=RBD Backup Tool

[Service]
Type=oneshot
ExecStart=/opt/bin/rbdbackuptool.sh
